# mdebian-iut-proxy/DEBIAN/postinst
#! /bin/sh

# Halt on error
set -e

# Disable the CNTLM service
systemctl stop cntlm
systemctl disable cntlm
